# docker/s6-assets/s6-rc.d/run-after-daemon-start/up
#!/command/execlineb -S0

foreground { s6-echo "Calling /etc/init/run-after-daemon-start" }
/etc/init/run-after-daemon-start.sh
